# supervisord.conf
[supervisord]
nodaemon=true

[program:octane-start]
command=php /var/www/html/artisan octane:start

